# pi-backup.timer
[Unit]
Description=Run Pi Backup daily

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
